# unix-daemon/resources/wifi-captive-daemon.service
[Unit]
Description=WiFi Captive Portal Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/wifi-captive-daemon
Restart=always
StandardOutput=append:/var/log/wifi-captive-daemon.log
StandardError=append:/var/log/wifi-captive-daemon.err

[Install]
WantedBy=multi-user.target
